# config/includes.chroot/etc/systemd/system/ocean-blue-tty.service
[Unit]
Description=Apply Ocean Blue color scheme to TTY consoles
After=systemd-vconsole-setup.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set-ocean-blue-tty.sh
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
